# assets/supervisor/kafka.ini
[program:kafka]
command=/usr/bin/start-kafka
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
